# build_files/users-setup/first-login.service
[Unit]
Description=First login setup script
# This condition ensures it only runs once per user
ConditionPathExists=!%h/.local/state/first-login-done

[Service]
Type=oneshot
# Run the script we created
ExecStart=/usr/bin/bash /usr/libexec/users-setup/first-login.sh
# Create a marker file so it never runs again for this user
ExecStartPost=/usr/bin/touch %h/.local/state/first-login-done

[Install]
WantedBy=default.target
